# systemd/zettanote-frontend.service
[Unit]
Description=ZettaNote Frontend Service
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/zettanote/current/frontend
ExecStart=/usr/bin/serve -s . -l 80
Restart=always
RestartSec=10
StandardOutput=append:/var/log/zettanote/frontend.log
StandardError=append:/var/log/zettanote/frontend-error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
